# ops/com.fitbit-plus.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Scheduled fitbit-plus sync: runs `python cli.py sync` every 6 hours.
  Install:   cp ops/com.fitbit-plus.sync.plist ~/Library/LaunchAgents/
             launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.fitbit-plus.sync.plist
  Remove:    launchctl bootout gui/$(id -u)/com.fitbit-plus.sync
  Logs:      ~/Library/Logs/fitbit-plus-sync.log
  Run now:   launchctl kickstart gui/$(id -u)/com.fitbit-plus.sync
  Exit 2 in the log means the 7-day token died - re-run `python cli.py auth`.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fitbit-plus.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/deekshith/Documents/Projects/agents/claude-agents/fitbit-plus/backend/.venv/bin/python</string>
        <string>cli.py</string>
        <string>sync</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/deekshith/Documents/Projects/agents/claude-agents/fitbit-plus/backend</string>
    <key>StartInterval</key>
    <integer>21600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/deekshith/Library/Logs/fitbit-plus-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/deekshith/Library/Logs/fitbit-plus-sync.log</string>
</dict>
</plist>
